# app/internal/supervisord/static/init.conf
[program:init]
command=mkdir -p /data/www/customer
directory=/
autostart=true
stderr_logfile=/dev/null
stdout_logfile=/dev/null
